# conf/celery-beat/celery_beat.conf
[supervisord]
nodaemon=true ;

[program:celery_beat]
user=dsuser ;
command=/home/dsuser/conf/celery-beat/celery_beat.sh ;
directory=/home/dsuser/ ;
priority=1 ;
