# conf/celery-worker/celery_worker.conf
[supervisord]
nodaemon=true ;

[program:celery_worker]
user=dsuser ;
command=/home/dsuser/conf/celery-worker/celery_worker.sh ;
directory=/home/dsuser/ ;
priority=1 ;