# conf/gunicorn/gunicorn.conf
[supervisord]
nodaemon=true ;

[program:gunicorn]
user=dsuser ;
command=/home/dsuser/conf/gunicorn/gunicorn.sh ;
directory=/home/dsuser/ ;
priority=1 ;
